# vpg/envs/scate-dockers/scate-backend/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:perf20web]
command=php /var/www/wmt-backend/queue_consumer.php gal